# docker/whitefield/supervisord.conf
[program:whitefield]
directory=/home/whitefield
command=pserve development.ini --reload
# user=whitefield
